# encoder/ioncube-encoder-api.service
[Unit]
Description=IonCube Encoder API Service
Documentation=https://github.com/your-repo/ioncube-encoder-api
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/ioncube-encoder-api
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ioncube-encoder-api
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=IONCUBE_PATH=/usr/local/bin/ioncube_encoder

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ioncube-encoder-api/uploads /opt/ioncube-encoder-api/temp /opt/ioncube-encoder-api/output
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
